# workspace/crontab/docker
* * * * * docker /usr/bin/php /var/www/blog/artisan schedule:run >> /dev/null 2>&1
